# crontab
# Run Symfony command every minute
* 8 * * * /usr/local/bin/php /app/bin/console app:my-command >> /var/log/cron.log 2>&1
