# crontab
# Run Symfony command every minute
* * * * * /usr/local/bin/php /app/bin/console app:my-command >> /var/log/cron.log 2>&1
